# crontab
# Backup databases daily at midnight
15 */12 * * * /root/cron_task.sh /root/db_backup.sh

